# scripts/autocode-headless.service
[Unit]
Description=Auto-Code Headless Development Service
After=network.target

[Service]
Type=simple
User=autocode
Group=autocode
WorkingDirectory=/opt/autocode/Auto-code-v1
Environment="ANTHROPIC_API_KEY=sk-ant-your-key-here"
Environment="AUTOCODE_REPO_ID=your-repository-id"
Environment="AUTOCODE_API_URL=http://localhost:8000"
Environment="AUTOCODE_USERNAME=admin"
Environment="AUTOCODE_PASSWORD=admin"
Environment="AUTOCODE_SLEEP_INTERVAL=300"
Environment="AUTOCODE_MAX_TICKETS=0"

ExecStart=/opt/autocode/Auto-code-v1/scripts/headless_dev.sh

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/autocode/headless.log
StandardError=append:/var/log/autocode/headless.error.log

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
